# mas/core/component/tools/svcs_tools/snmpagent.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    Copyright 2005 Mobeon AB.  All rights reserved.
    Use is subject to license terms.
    mas.xml : MAS manifest, Hakan Tunell 
    updated: 2006-01-26
-->

<service_bundle type='manifest' name='MOBYmas:snmpagent'>

   <service
        name='application/mobeon/MOBYmas'
        type='service'
        version='1'>

   <single_instance />

   <dependency
      name='filesystem'
      grouping='require_all'
      restart_on='none'
      type='service'>
      <service_fmri value='svc:/system/filesystem/local' />
   </dependency>

   <exec_method 
      type='method' 
      name='start'
      exec='/etc/sfw/mobeon/MOBYmas/snmpagent start'
      timeout_seconds='120' />

   <exec_method 
      type='method' 
      name='stop'
      exec=':kill'
      timeout_seconds='120' />

   <instance name='snmpagent' enabled='false' />

   <stability value='Unstable' />

   <template>
      <common_name>
         <loctext xml:lang='C'>M3 Media Access Server/SNMP agent </loctext>
      </common_name>
      <documentation>
         <manpage title='CPI document: MAS Feature ' section=' ' manpath=' ' />
      </documentation>
   </template>

</service>
</service_bundle>
